# supervisord.conf
[supervisord]
nodaemon=true

[program:shibd]
command=service shibd start && tail -f /var/log/shibboleth/shibd.log

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
